# scripts/.crontab
# Awesome barberini crontab file. See setup_cron.sh for installing this.
*/5 * * * * /root/bp-barberini/scripts/pull_repo.sh
0 3 * * * /root/bp-barberini/scripts/cron.sh daily
0 * * * * /root/bp-barberini/scripts/cron.sh hourly
